# nat123/nat123.service
[Unit]
Description=Nat123
After=network.target

[Service]
ExecStart=/usr/local/bin/nat123

[Install]
WantedBy=multi-user.target

